# twingate-client/rootfs/etc/services.d/twingate-client/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

SERVICE_KEY=$(bashio::config 'service_key')
LOG_LEVEL=$(bashio::config 'log_level')
CONFIGURE_DNS=$(bashio::config 'configure_hassio_dns')
export LOG_LEVEL

if [ -z "${SERVICE_KEY}" ]; then
    bashio::log.fatal "No service_key configured. Generate one in Twingate Admin > Service Accounts."
    exit 1
fi

bashio::log.info "Creating TUN device..."
mkdir -p /dev/net
if [ ! -c /dev/net/tun ]; then
    mknod /dev/net/tun c 10 200
    chmod 600 /dev/net/tun
fi

# Ensure the container hostname is resolvable (Twingate requires this)
if ! getent hosts "$(hostname)" > /dev/null 2>&1; then
    echo "127.0.0.1 $(hostname)" >> /etc/hosts
fi

# Override DNS to avoid Twingate detecting Docker/HA stub resolvers
# (127.0.0.11 or 172.30.x.x) as its own DNS stubs, which prevents
# the tunnel interface from coming up.
bashio::log.info "Configuring DNS for Twingate..."
cat > /etc/resolv.conf << 'DNSEOF'
nameserver 1.1.1.1
nameserver 8.8.8.8
DNSEOF

bashio::log.info "Writing service key..."
echo "${SERVICE_KEY}" > /etc/twingate/service_key.json
chmod 600 /etc/twingate/service_key.json

bashio::log.info "Configuring Twingate headless client..."
twingate setup --headless /etc/twingate/service_key.json

# Stop any previous daemon instance to clear stale state before starting.
# Without this, socket/lock files from a crashed daemon prevent twingate start
# from succeeding on s6-triggered restarts.
bashio::log.info "Stopping any previous daemon instance..."
twingate stop 2>/dev/null || true
sleep 2

bashio::log.info "Starting Twingate client..."
twingate start

# twingate start daemonises and exits. Keep the s6 service alive
# by monitoring the daemon using twingate status.
sleep 15
bashio::log.info "Initial status: $(twingate status 2>&1)"

# --------------------------------------------------------------------------
# Configure hassio-dns to use Twingate's internal DNS resolvers.
# Other add-ons (e.g. OpenClaw) use hassio-dns (172.30.32.3) which does not
# go through Twingate's DNS interception. Adding Twingate's DNS IPs as
# upstream servers in hassio-dns ensures all containers can resolve
# Twingate-protected resource domains to their virtual tunnel IPs.
# --------------------------------------------------------------------------
if [ "${CONFIGURE_DNS}" = "true" ]; then
    configure_hassio_dns() {
        # Extract Twingate internal DNS IPs from the sdwan route table.
        # These are the 100.95.0.x host routes that Twingate creates.
        local tg_dns_ips
        tg_dns_ips=$(ip route show | grep -oP '100\.95\.0\.\d+(?= dev sdwan)' | head -4)

        if [ -z "${tg_dns_ips}" ]; then
            bashio::log.warning "No Twingate DNS IPs found in route table. Skipping hassio-dns configuration."
            return 1
        fi

        # Save current hassio-dns server config for restoration on shutdown
        local current_servers
        current_servers=$(curl -s -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
            http://supervisor/dns/info \
            | sed -n 's/.*"servers":\[\([^]]*\)\].*/\1/p' 2>/dev/null)
        echo "${current_servers}" > /tmp/hassio_dns_original_servers

        # Build the API payload with Twingate DNS servers
        local dns_payload="["
        local first=true
        for ip in ${tg_dns_ips}; do
            if [ "${first}" = "true" ]; then
                first=false
            else
                dns_payload+=","
            fi
            dns_payload+="\"dns://${ip}\""
        done
        dns_payload+="]"

        bashio::log.info "Configuring hassio-dns with Twingate resolvers: ${dns_payload}"

        local result
        result=$(curl -s -o /dev/null -w '%{http_code}' \
            -X POST \
            -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
            -H "Content-Type: application/json" \
            -d "{\"servers\": ${dns_payload}}" \
            http://supervisor/dns/options)

        if [ "${result}" = "200" ]; then
            bashio::log.info "Hassio-dns configured with Twingate DNS resolvers"
            # Restart hassio-dns to apply the new config
            curl -s -X POST \
                -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
                http://supervisor/dns/restart > /dev/null 2>&1
            return 0
        else
            bashio::log.warning "Failed to configure hassio-dns (HTTP ${result})"
            return 1
        fi
    }

    configure_hassio_dns
fi

while true; do
    STATUS=$(twingate status 2>&1)
    if [ "$STATUS" = "not-running" ]; then
        bashio::log.error "Twingate daemon is not running"
        bashio::log.warning "Daemon log tail:"
        tail -20 /var/log/twingated.log 2>/dev/null || true
        exit 1
    fi
    sleep 30
done
