# twingate-client/rootfs/etc/services.d/twingate-client/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

SERVICE_KEY=$(bashio::config 'service_key')
LOG_LEVEL=$(bashio::config 'log_level')
export LOG_LEVEL

if [ -z "${SERVICE_KEY}" ]; then
    bashio::log.fatal "No service_key configured. Generate one in Twingate Admin > Service Accounts."
    exit 1
fi

bashio::log.info "Creating TUN device..."
mkdir -p /dev/net
if [ ! -c /dev/net/tun ]; then
    mknod /dev/net/tun c 10 200
    chmod 600 /dev/net/tun
fi

# Ensure the container hostname is resolvable (Twingate requires this)
if ! getent hosts "$(hostname)" > /dev/null 2>&1; then
    echo "127.0.0.1 $(hostname)" >> /etc/hosts
fi

# Override DNS to avoid Twingate detecting Docker/HA stub resolvers
# (127.0.0.11 or 172.30.x.x) as its own DNS stubs, which prevents
# the tunnel interface from coming up.
bashio::log.info "Configuring DNS for Twingate..."
cat > /etc/resolv.conf << 'DNSEOF'
nameserver 1.1.1.1
nameserver 8.8.8.8
DNSEOF

bashio::log.info "Writing service key..."
echo "${SERVICE_KEY}" > /etc/twingate/service_key.json
chmod 600 /etc/twingate/service_key.json

bashio::log.info "Configuring Twingate headless client..."
twingate setup --headless /etc/twingate/service_key.json

# Stop any previous daemon instance to clear stale state before starting.
# Without this, socket/lock files from a crashed daemon prevent twingate start
# from succeeding on s6-triggered restarts.
bashio::log.info "Stopping any previous daemon instance..."
twingate stop 2>/dev/null || true
sleep 2

bashio::log.info "Starting Twingate client..."
twingate start

# twingate start daemonises and exits. Keep the s6 service alive
# by monitoring the daemon using twingate status.
sleep 15
bashio::log.info "Initial status: $(twingate status 2>&1)"

while true; do
    STATUS=$(twingate status 2>&1)
    if [ "$STATUS" = "not-running" ]; then
        bashio::log.error "Twingate daemon is not running"
        bashio::log.warning "Daemon log tail:"
        tail -20 /var/log/twingated.log 2>/dev/null || true
        exit 1
    fi
    sleep 30
done
